# layers/meta-cclrte/recipes-realtime/rt-verify/files/rt-verify.service
[Unit]
Description=CCLRTE RT Latency Verification (cyclictest)
Documentation=https://github.com/user/codesys-control-linux-rte/docs/LIMITATIONS.md
After=rt-setup.service
# No ConditionPathExists — WebUI deletes the result file before starting this
# so it can be re-triggered on demand at any time.

[Service]
Type=oneshot
ExecStart=/usr/sbin/run-cyclictest.sh
StandardOutput=journal
StandardError=journal
# 3 phases × 60s + 30s margin
TimeoutStartSec=220

[Install]
WantedBy=multi-user.target
