# layers/meta-cclrte/recipes-codesys/codesys-control/files/services/codesyscontrol.service
[Unit]
Description=CODESYS Control for Linux SL Runtime
Documentation=https://store.codesys.com
After=network.target rt-setup.service ethercat.service
Wants=ethercat.service
Requires=rt-setup.service
# Do not start if binary is not yet installed.
# install-codesys-runtime.sh creates codesyscontrol symlink → codesyscontrol.bin
ConditionPathExists=/opt/codesys/bin/codesyscontrol

[Service]
# codesyscontrol.bin does not self-daemonize — the init.d script runs it with &
# Use Type=simple so systemd tracks the process directly.
Type=simple
WorkingDirectory=/var/opt/codesys
RuntimeDirectory=codesys
RuntimeDirectoryMode=0755

ExecStartPre=/usr/sbin/codesys-setup.sh
# Config file is the first positional argument (no -d / -c flags)
ExecStart=/opt/codesys/bin/codesyscontrol /etc/codesyscontrol/CODESYSControl.cfg
EnvironmentFile=-/etc/default/codesyscontrol
Environment=LD_LIBRARY_PATH=/opt/codesys/lib

Restart=on-failure
RestartSec=5s

# Real-time scheduling — CPU3, SCHED_FIFO 80
# Reinforced by drop-in: /etc/systemd/system/codesyscontrol.service.d/rt-override.conf
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80
CPUAffinity=3
LimitMEMLOCK=infinity
LimitNOFILE=65536
LimitSTACK=8388608

[Install]
WantedBy=multi-user.target
