# layers/meta-cclrte/recipes-codesys/codesys-control/files/codesyscontrol.service
[Unit]
Description=CODESYS Control for Linux SL Runtime
Documentation=https://store.codesys.com
After=network.target rt-setup.service ethercat.service codesysgateway.service
Wants=ethercat.service codesysgateway.service
Requires=rt-setup.service
# Do not start if binary is not yet installed by CODESYS IDE
ConditionPathExists=/opt/codesys/bin/codesyscontrol

[Service]
Type=forking
PIDFile=/run/codesys/codesyscontrol.pid
RuntimeDirectory=codesys
RuntimeDirectoryMode=0755

ExecStartPre=/usr/sbin/codesys-setup.sh
ExecStart=/opt/codesys/bin/codesyscontrol -d -c /etc/CODESYSControl.cfg
ExecStop=/opt/codesys/bin/codesyscontrol -k

Restart=on-failure
RestartSec=5s

# Real-time scheduling — CPU3, SCHED_FIFO 80
# Reinforced by drop-in: /etc/systemd/system/codesyscontrol.service.d/rt-override.conf
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80
CPUAffinity=3
LimitMEMLOCK=infinity
LimitNOFILE=65536
LimitSTACK=8388608

WatchdogSec=10s
NotifyAccess=main

[Install]
WantedBy=multi-user.target
